# config/systemd/user/autostart.service
[Unit]
Description=Autostart
After=home.mount

[Service]
ExecStart=/bin/bash autostart.sh

[Install]
WantedBy=default.target